# machineConfigs/servers/k3s-server/etc/systemd/system/wol.service
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /usr/sbin/ethtool --change enp1s0 wol g

[Install]
WantedBy=basic.target
