# machineConfigs/desktop/home/arthur/.config/autostart/com.github.wwmm.easyeffects.desktop
[Desktop Entry]
Type=Application
Name=com.github.wwmm.easyeffects
Exec=flatpak run --command=easyeffects com.github.wwmm.easyeffects --gapplication-service
X-Flatpak=com.github.wwmm.easyeffects
