# machineConfigs/k3s/var/spool/cron/crontabs/root
@reboot     sleep 10;   /bin/bash -c 'sh /home/arthur/discord/discord_bot.sh booted'

